# upstart/pyfd.conf
description "Pyfortunes daemon"

start on runlevel [2345]
stop  on runlevel [016] 

expect fork

exec /usr/bin/pyfd /etc/pyfd.conf 2> /var/log/pyfd.log
